# packages/images/mkosi.profiles/azure/mkosi.postinst
#!/bin/bash

mkosi-chroot systemctl enable "azure-provisioning.service"
ln -sf "/etc/systemd/system/azure-provisioning.service" "$BUILDROOT/etc/systemd/system/minimal.target.wants/"
