# packages/images/tdx-dummy/kettle-log-dumper.service
[Unit]
Description=Kettle Log Dumper (Debug)
After=kettle-launcher.service

[Service]
Type=simple
ExecStart=/usr/local/bin/dump-kettle-logs.sh
StandardOutput=journal+console
StandardError=journal+console

[Install]
WantedBy=minimal.target
