# packages/images/tdx-dummy/kettle-launcher.service
[Unit]
Description=Kettle Application Launcher
After=network-setup.service
Wants=network-setup.service

[Service]
Type=exec
User=root
Group=root
ExecStart=/usr/bin/kettle launch /usr/lib/kettle/manifest.json --port 3001 --db-dir /var/lib/kettle/db
Restart=on-failure
RestartSec=10
StandardOutput=append:/var/log/kettle.log
StandardError=append:/var/log/kettle-error.log

[Install]
WantedBy=minimal.target
