# packages/images/base/mkosi.skeleton/etc/systemd/system/persistent-mount.service
[Unit]
Description=Mount Persistent Data Disk
DefaultDependencies=no
Conflicts=shutdown.target
Before=minimal.target
After=local-fs-pre.target local-fs.target

[Service]
Type=oneshot
ExecStart=/usr/bin/mount-persistent
RemainAfterExit=yes
StandardOutput=journal+console
StandardError=journal+console

[Install]
WantedBy=minimal.target
