# packages/images/mkosi.profiles/azure-devtools/mkosi.skeleton/etc/systemd/system/ssh-keygen.service
[Unit]
Description=Generate SSH host keys on first boot
DefaultDependencies=no
Before=ssh.service sshd.service
ConditionPathExists=!/etc/ssh/ssh_host_ed25519_key

[Service]
Type=oneshot
ExecStart=/usr/bin/ssh-keygen -A
RemainAfterExit=yes
StandardOutput=append:/var/log/ssh-keygen.log
StandardError=append:/var/log/ssh-keygen.log

[Install]
WantedBy=multi-user.target
