# packages/images/tdx-kettle/cloud-launcher.service
[Unit]
Description=Cloud Configuration Service (Fetch Manifest from Metadata)
After=network-setup.service
Wants=network-setup.service
Before=kettle-launcher.service

[Service]
Type=oneshot
User=root
Group=root
ExecStart=/usr/bin/cloud-launcher
RemainAfterExit=yes
StandardOutput=append:/var/log/cloud-launcher.log

[Install]
WantedBy=minimal.target
