# packages/images/tdx-kettle/kettle-launcher.service
[Unit]
Description=Kettle Application Launcher
After=network-setup.service cloud-config.service
Wants=network-setup.service cloud-config.service

[Service]
Type=exec
User=root
Group=root
# Load cloud-specific environment variables (MANIFEST) if available
EnvironmentFile=-/etc/kettle/cloud-config.env
ExecStart=/usr/bin/kettle-init
Restart=on-failure
RestartSec=10
StandardOutput=append:/var/log/kettle.log
StandardError=append:/var/log/kettle-error.log

[Install]
WantedBy=minimal.target
